# output/run/systemd/generator/rpc_pipefs.target
# Automatically generated by rpc-pipefs-generator

[Unit]
Requires=run-rpc_pipefs.mount
After=run-rpc_pipefs.mount
